# root/etc/services.d/caddy/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /usr/local/bin/caddy -agree \
	-conf /config/Caddyfile && \
	php-fpm7

